# deb/install/postinst
#!/bin/bash

case "${1}" in
  configure)
    echo "Enabling waggle-firewall service"
    systemctl enable waggle-firewall.service
    ;;
esac
