# deb/install/prerm
#!/bin/bash

systemctl stop waggle-firewall.service || true
systemctl disable waggle-firewall.service || true
